# deploy-example/test.conf
[program:test]
command=python /root/test/test.py
stderr_logfile = /root/test/test-stderr.log
stdout_logfile = /root/test/test-stdout.log
